# dist/gtm/bin/gstartup
# chkconfig: 235 96 10
# description: Start or stop the GTM 
#
### BEGIN INIT INFO
# Default-Start: 2 3 5
# Default-Stop: 0 1 6
### END INIT INFO
#

if [ "x$1" = "xstart" ]
then /gtm/bin/gstart
fi
if [ "x$1" = "xstop" ]
then /gtm/bin/gstop
fi

exit 0